# examples/systemd/orange-xp-daily.timer
[Unit]
Description=Schedule Orange XP OpenClaw workflow

[Timer]
OnCalendar=*-*-* 09:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
